# crontab
17 * * * * /usr/lib/emax_em3371_decoder/set_device_to_client.sh
